# ledger/LEDGER-0024-dual-push-drift-prevention/artifacts/yousirjuan-dual-push-sync.timer
[Unit]
Description=marvelous-empire dual-push drift detector — runs every 5 min

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
RandomizedDelaySec=30s
Persistent=true

[Install]
WantedBy=timers.target
