# ledger/LEDGER-0024-dual-push-drift-prevention/artifacts/yousirjuan-dual-push-sync.service
[Unit]
Description=marvelous-empire dual-push drift detector + auto-mirror
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/yousirjuan-sync/sync-and-drift.sh
TimeoutStartSec=600
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
