# ledger/LEDGER-0009-watchdog-tamer/artifacts/com.yousirjuan.watchdog-tamer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.yousirjuan.watchdog-tamer</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__PLAYBOOK_PATH__</string>
  </array>

  <key>StartInterval</key>
  <integer>900</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>__HOME__</string>
    <key>OLLAMA_URL</key>
    <string>http://127.0.0.1:11434</string>
    <key>TAMER_MODEL</key>
    <string>llama3.2:3b</string>
  </dict>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/yousirjuan-watchdog-tamer.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/yousirjuan-watchdog-tamer.stderr.log</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
