# ledger/LEDGER-0008-watchdog-control-surface/artifacts/com.yousirjuan.watchdog-state-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.yousirjuan.watchdog-state-server</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__PLAYBOOK_PATH__</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>__HOME__</string>
    <key>WATCHDOG_BIND</key>
    <string>0.0.0.0</string>
    <key>WATCHDOG_PORT</key>
    <string>9876</string>
  </dict>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/yousirjuan-watchdog-server.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/yousirjuan-watchdog-server.stderr.log</string>

  <key>ProcessType</key>
  <string>Background</string>

  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
